# deploy/bench-central.service
[Unit]
Description=Bench Race Central - Benchmark coordinator and web UI
Documentation=https://github.com/todddb/bench-race
After=network.target

[Service]
Type=simple
User=bench
Group=bench
WorkingDirectory=/opt/bench-race
ExecStart=/opt/bench-race/bin/control central start --foreground
ExecStop=/opt/bench-race/bin/control central stop
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Environment variables
Environment=CENTRAL_HOST=0.0.0.0
Environment=CENTRAL_PORT=8080
# Optional: SERVICE_CONTROL_TOKEN for API authentication
# Environment=SERVICE_CONTROL_TOKEN=changeme

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/bench-race/run /opt/bench-race/logs /opt/bench-race/results

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bench-central

[Install]
WantedBy=multi-user.target
