# deploy/com.bench-race.vllm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bench-race.vllm</string>

    <key>ProgramArguments</key>
    <array>
        <string>__VENV_PATH__/bin/python</string>
        <string>-m</string>
        <string>vllm.entrypoints.openai.api_server</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>__VLLM_PORT__</string>
        <string>--download-dir</string>
        <string>__MODEL_DIR__</string>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HF_HOME</key>
        <string>__MODEL_DIR__</string>
        <key>PATH</key>
        <string>__VENV_PATH__/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>__VENV_PATH__</string>

    <key>StandardOutPath</key>
    <string>/tmp/bench-race-vllm.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/bench-race-vllm.stderr.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
